# dist/ebpfsentinel-token.service.d/override.conf
# Systemd drop-in override for BPF token-based loading.
#
# Installs at /etc/systemd/system/ebpfsentinel.service.d/bpf-token.conf.
#
# Runs the token setup script as a privileged pre-exec, then strips
# the service down to the minimal capability set. The token fd is
# inherited through the bpffs mount that the setup script pinned, so
# the main process no longer needs CAP_BPF / CAP_NET_ADMIN.
#
# Apply with:
#   sudo mkdir -p /etc/systemd/system/ebpfsentinel.service.d
#   sudo cp override.conf /etc/systemd/system/ebpfsentinel.service.d/bpf-token.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart ebpfsentinel
#
# Verify with:
#   systemctl status ebpfsentinel
#   curl -s localhost:9090/metrics | grep bpf_token_used

[Service]
# Run token setup before the main agent. Uses `-` prefix to tolerate
# the idempotent "already mounted" exit.
ExecStartPre=/usr/local/bin/ebpfsentinel-token-setup.sh /sys/fs/bpf/ebpfsentinel

# Strip ambient capabilities — the token fd replaces CAP_BPF +
# CAP_NET_ADMIN. CAP_SYS_ADMIN is only needed by the pre-exec script.
# CAP_NET_RAW stays for pcap-based manual packet capture.
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW

# Limit the process to the delegated bpffs mount
ReadOnlyPaths=/sys/fs/bpf/ebpfsentinel

# Stronger sandbox on top of the base service hardening
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
